# rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv sh

env HOME=/srv s6-setuidgid srv /usr/bin/caddy run --config /srv/caddy/Caddyfile
