# rootfs/etc/services.d/xray/run
#!/usr/bin/with-contenv sh

sleep 5

CERTFILE="/srv/xray/fullchain.crt"

if [ -f $CERTFILE ]; then
    s6-setuidgid srv /usr/bin/xray -config /srv/xray/config.json
fi
